# internal/podman/quadlets/lerd-php-fpm.container.tmpl
[Unit]
Description=Lerd PHP {{.Version}} FPM
After=network.target
# Tolerate transient restart cascades (e.g. EnsurePathMounted called rapidly
# with rotating paths) without permanently killing FPM and the worker units
# that BindsTo it. systemd's default 5-restarts-in-10s is way too tight; this
# allows up to 20 restarts in a 5-minute window, which absorbs realistic
# cascades but still gives up on a genuinely broken FPM (corrupt image,
# bad ini) instead of looping forever and burning CPU.
StartLimitIntervalSec=300
StartLimitBurst=20

[Container]
Image=lerd-php{{.VersionShort}}-fpm:local
ContainerName=lerd-php{{.VersionShort}}-fpm
{{.HostNameLine}}
Network=lerd
Volume=%h/.local/share/lerd/hosts:/etc/hosts:ro,z
Volume=%h:%h:rw
Volume={{.XdebugIniPath}}:/usr/local/etc/php/conf.d/99-xdebug.ini:ro
Volume={{.UserIniPath}}:/usr/local/etc/php/conf.d/98-lerd-user.ini:ro
Volume={{.DumpsDir}}:/usr/local/etc/lerd:ro
Volume={{.DumpsIniPath}}:/usr/local/etc/php/conf.d/97-lerd-dump.ini:ro
Volume={{.DevtoolsIniPath}}:/usr/local/etc/php/conf.d/96-lerd-devtools.ini:ro
Volume={{.SpxIniPath}}:/usr/local/etc/php/conf.d/zz-lerd-spx.ini:ro
Volume={{.SpxDataDir}}:/var/spx:rw
Volume={{.ZshHistoryDir}}:/root/.zsh_state:rw
# Persistent home for an optional in-container musl bun (lerd php:bun install).
# bun installs and self-updates here; survives image rebuilds and is shared
# across PHP versions. Empty until the user opts in.
Volume={{.BunVolumeDir}}:/root/.bun:rw
# Persistent Playwright browser cache for opt-in Pest browser testing (lerd
# pest:browser install). Holds the downloaded registry and lerd's musl-chromium
# shims so they survive image rebuilds. Empty until the user opts in.
Volume={{.PlaywrightVolumeDir}}:/root/.cache/ms-playwright:rw
Environment=HOME=/root
PodmanArgs=--security-opt=label=disable
Exec=php-fpm -F -R

[Service]
Restart=always
# 2-second backoff between restarts so a tight crash loop can't burn CPU.
RestartSec=2

[Install]
WantedBy=default.target
